# deploy/systemd/gpu-worker.service
[Unit]
Description=GPU Worker - WhisperX transcription and CoPE-A classification
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/pipeline
EnvironmentFile=/opt/pipeline/.env
ExecStart=/opt/pipeline/venv/bin/python -m src.gpu_worker
Restart=always
RestartSec=30

# GPU workers need time to load models
TimeoutStartSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gpu-worker

[Install]
WantedBy=multi-user.target
