# deploy/systemd/transfer-worker.service
[Unit]
Description=Transfer Worker - Downloads archives from AWS EC2
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/pipeline
EnvironmentFile=/opt/pipeline/.env
ExecStart=/opt/pipeline/venv/bin/python -m src.transfer_sounds
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=transfer-worker

[Install]
WantedBy=multi-user.target
