# deploy/systemd/unpack-worker.service
[Unit]
Description=Unpack Worker - Extracts archives and converts MP3 to Opus
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/pipeline
EnvironmentFile=/opt/pipeline/.env
ExecStart=/opt/pipeline/venv/bin/python -m src.unpack_worker
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=unpack-worker

[Install]
WantedBy=multi-user.target
