# ansible/roles/gpu_worker/templates/gpu-worker.service.j2
[Unit]
Description=Audio Pipeline GPU Worker (WhisperX + CoPE-A)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ pipeline_dir }}

# Environment variables
Environment=PYTHONUNBUFFERED=1
Environment=CUDA_VISIBLE_DEVICES={{ cuda_visible_devices }}

# Redis connection
Environment=REDIS_HOST={{ redis_host }}
Environment=REDIS_PORT={{ redis_port }}

# PostgreSQL connection
Environment=POSTGRES_HOST={{ postgres_host }}
Environment=POSTGRES_PORT={{ postgres_port }}
Environment=POSTGRES_DB={{ postgres_db }}
Environment=POSTGRES_USER={{ postgres_user }}
Environment=POSTGRES_PASSWORD={{ postgres_password }}

# S3 configuration
Environment=S3_ENDPOINT={{ s3_endpoint }}
Environment=S3_ACCESS_KEY={{ s3_access_key }}
Environment=S3_SECRET_KEY={{ s3_secret_key }}
Environment=S3_BUCKET={{ s3_bucket }}

# Local storage paths
Environment=SCRATCH_ROOT={{ scratch_root }}
Environment=MODELS_ROOT={{ models_root }}

# Processing settings
Environment=WHISPERX_MODEL={{ whisperx_model }}
Environment=GPU_BATCH_SIZE={{ gpu_batch_size }}
Environment=FFMPEG_WORKERS={{ ffmpeg_workers }}
Environment=OPUS_BITRATE={{ opus_bitrate }}
Environment=COPE_MODEL={{ gemma_model_path }}
Environment=COPE_ADAPTER={{ cope_adapter_path }}

# Logging
Environment=LOG_DIR={{ log_dir }}

# Execute
ExecStart={{ pipeline_venv }}/bin/python -m src.gpu_worker

# Restart policy
Restart=always
RestartSec=30

# Model loading can take time
TimeoutStartSec=300

# Resource limits
LimitNOFILE=65535

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gpu-worker

[Install]
WantedBy=multi-user.target
